# bench-runner.service
[Unit]
Description=Autonomous Benchmark Runner
# Wait for network before starting so the metadata API is reachable.
# cloud-init.target is intentionally omitted — it has After=multi-user.target
# which creates an ordering cycle that causes systemd to silently drop this
# service. The runner script retries the metadata fetch for up to 2 minutes,
# which is sufficient to wait out any remaining cloud-init activity.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/run_benchmarks.sh
# Don't restart on failure — avoids runaway billing if the script has a bug.
Restart=no
RemainAfterExit=yes
# 6 hours: generous upper bound for even very slow models; safety-net 'at' job
# fires at 5 hours regardless.
TimeoutStartSec=21600
# Load VULTR_API_KEY and any other env vars written to /etc/environment
EnvironmentFile=-/etc/environment
# Write service output to the journal (viewable with journalctl -u bench-runner)
# The script also tees to /var/log/bench-runner.log for easy SSH inspection.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bench-runner

[Install]
WantedBy=multi-user.target
